# task2/supervisord.conf
[supervisord]
nodaemon=true

[program:docker]
command=dockerd -D

[program:apparmor]
command=/start.sh

